# backend/supervisor_django.conf
[program:django-backend]
; The program directive specifies the command to run
command=/Users/raghavjasthi/Desktop/Fusion-1/venv/bin/python /Users/raghavjasthi/Desktop/Fusion-1/manage.py runserver 0.0.0.0:8000

; The directory where the process will run
directory=/Users/raghavjasthi/Desktop/Fusion-1

; The user that runs the program
user=raghavjasthi

; The number of seconds the program needs to stay running before supervisord considers it a success
startsecs=10

; Whether supervisord will automatically start this program when supervisord starts
autostart=true

; Whether supervisord will automatically restart the program if it exits unexpectedly
autorestart=true

; The maximum number of seconds the program can run
; Set to 0 for unlimited runtime
stopasgroup=true

; Output logs
stdout_logfile=/Users/raghavjasthi/Desktop/Fusion-1/logs/django_stdout.log
stderr_logfile=/Users/raghavjasthi/Desktop/Fusion-1/logs/django_stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

; The SIGTERM signal is sent to the program when a stop is requested
stopsignal=QUIT

; Whether the process will be restarted even after intentional stops
priority=999
